# datalogger.service
[Unit]
Description=DataLogger Temperature Monitoring Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/DataLogger
Environment="PATH=/home/pi/DataLogger/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/pi/DataLogger/.venv/bin/python /home/pi/DataLogger/data-logger-project/app_enhanced.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target